# config/autostart/xfce4-panel.desktop
[Desktop Entry]
Type=Application
Name=Xfce Panel
Exec=xfce4-panel
OnlyShowIn=XFCE;
RunHook=0
