# config/autostart/xfdesktop.desktop
[Desktop Entry]
Type=Application
Name=Xfce Desktop
Exec=xfdesktop
OnlyShowIn=XFCE;
RunHook=0
