# config/autostart/xfsettingsd.desktop
[Desktop Entry]
Type=Application
Name=Xfce Settings Daemon
Exec=xfsettingsd
OnlyShowIn=XFCE;
RunHook=0
